# packaging/rpm/SOURCES/rhelmon.service
[Unit]
Description=rhelmon - RHEL system monitoring agent
Documentation=https://github.com/amit25sep/rhelmon
After=network.target

[Service]
Type=simple
User=rhelmon
Group=rhelmon
ExecStart=/usr/bin/rhelmon-start
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rhelmon
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadOnlyPaths=/proc /sys
ReadWritePaths=/var/log/rhelmon /etc/rhelmon

[Install]
WantedBy=multi-user.target
